# com.gensyn.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.gensyn.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>sleep 30 && open -a Terminal /Users/satbaba2025b/autostart_setup/gensyn_watchdog.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
